# GW.Sqlbot/sqlbot.service
[Unit]
Description=SQLBot - 智能数据查询系统
Documentation=https://github.com/your-org/sqlbot
After=network.target
Wants=network.target

[Service]
Type=simple
User=sqlbot
Group=sqlbot
WorkingDirectory=/opt/sqlbot
Environment=PATH=/opt/sqlbot/venv/bin
ExecStart=/opt/sqlbot/venv/bin/python /opt/sqlbot/start_sqlbot.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sqlbot

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sqlbot/logs /opt/sqlbot/data /opt/sqlbot/backup
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target